# contrib/plists/com.geoffjay.agentd-ask.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.geoffjay.agentd-ask</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Agent.app/Contents/MacOS/agentd-ask</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/agentd-ask.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/agentd-ask.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>PORT</key>
        <string>7001</string>
        <key>NOTIFY_SERVICE_URL</key>
        <string>http://localhost:7004</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/usr/local</string>
</dict>
</plist>
